# deploy/pursuit.service
[Unit]
Description=Web service for hosting of PureScript API documentation

[Service]
Type=simple
User=www-data
# The server has 2 vCPUs and 4GB of RAM. The heap limit needs to leave room
# for nginx and the rest of the system, so that the RTS runs out of heap
# before the kernel's OOM killer steps in.
ExecStart=/usr/local/bin/pursuit +RTS -N2 -A64m -M3G -RTS
Restart=always
RestartSec=5s
Environment="PURSUIT_APPROOT=https://pursuit.purescript.org"
Environment="PURSUIT_HOST=127.0.0.1"
Environment="PURSUIT_PORT=3000"
Environment="PURSUIT_DATA_DIR=/var/www/pursuit/data"
Environment="PURSUIT_MINIMUM_COMPILER_VERSION=0.10.0"
EnvironmentFile=/var/www/pursuit/secrets

[Install]
WantedBy=multi-user.target
